# bdprochot.service
[Unit]
Description=bdprochot

[Service]
ExecStart=/etc/bdprochot.sh

[Install]
WantedBy=multi-user.target
